# root/usr/service/mongrel/run
#!/bin/sh

exec 2>&1

exec mongrel_rails start \
  --chdir /giggi \
  --port 80 \
  --log /var/giggi/log/production.log \
  --environment production \
  --user webapp \
  --group giggi
