# root/etc/inittab
::sysinit:/etc/rcS
::wait:/giggi/script/startup.rb
::respawn:/bin/runsvdir /var/service
::askfirst:/bin/sh
::shutdown:/giggi/script/shutdown.rb
::shutdown:/bin/umount -a -r
::restart:/sbin/init
