# root/usr/service/telnetd/run
#!/bin/sh

exec 2>&1

echo "Starting telnetd"
exec /sbin/telnetd -F
